# systemd/vine-watcher.service
[Unit]
Description=Vine Watcher
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/vine-watcher-telegram
ExecStart=/opt/vine-watcher-telegram/scripts/run-service.sh
Restart=on-failure
RestartSec=20
User=vinewatcher
Environment=NODE_ENV=production
Environment=HOME=/var/lib/vinewatcher
Environment=HEADLESS=false
Environment=VINE_WATCHER_XVFB=true

[Install]
WantedBy=multi-user.target
