# systemd/vine-watcher.service
[Unit]
Description=Vine Watcher
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/vine-watcher-telegram
ExecStart=/opt/vine-watcher-telegram/scripts/run-service.sh
Restart=on-failure
RestartSec=20
User=vinewatcher
Environment=NODE_ENV=production
Environment=HOME=/var/lib/vinewatcher
Environment=HEADLESS=false
Environment=VINE_WATCHER_XVFB=true
Environment=APP_DIR=/opt/vine-watcher-telegram
KillMode=control-group
TimeoutStopSec=45
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/vine-watcher-telegram/data /var/lib/vinewatcher
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
